# services/victoriametrics-stack.target
[Unit]
Description=VictoriaMetrics single-host stack



[Install]
WantedBy=multi-user.target

