# services/vm-one_week.service
[Unit]
Description=VictoriaMetrics one week retention (raw)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=victoriametrics
Group=victoriametrics
ExecStart=/opt/victoria-metrics-prod \
  -storageDataPath=/srv/victoria/one_week \
  -retentionPeriod=1w \
  -httpListenAddr=:8428
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
